# units/xfce-init.target
[Unit]
Description=Set up session and settings manager for Xfce
Before=xfce.target

[Install]
WantedBy=xfce.target
