# units/xfce.target
[Unit]
Description=Xfce Desktop Environment
After=xfce-init.target
