# units/xfce-pre.target
[Unit]
Description=Set up plumbing for Xfce
Before=xfce-init.target

[Install]
WantedBy=xfce-init.target
