# units/xfce.target
[Unit]
Description=Xfce Desktop Environment
