# units/xfce-init.target
[Unit]
Description=Set up session and settings manager for Xfce
After=xfce4_session_set_SESSIONMANAGER_var_service

[Install]
WantedBy=xfce.target
