# units/xfce-pre.target
[Unit]
Description=Set up plumbing for Xfce
Before=xfce.target
